# etc/profile.d/rust-env.sh
export RUSTUP_HOME=/media/data/local/rustup
export CARGO_HOME=/media/data/local/cargo
export RUSTUP_DIST_SERVER="https://rsproxy.cn"
export RUSTUP_UPDATE_ROOT="https://rsproxy.cn/rustup"
source "$CARGO_HOME/env"

